# chef/chef.conffiles
etc/chef/client.rb
etc/chef/solo.rb
etc/init.d/chef-client
