# chef/chef-server.conffiles
etc/chef/server.rb
etc/chef/indexer.rb
etc/init.d/chef-indexer.rb
